# plugins/realtime/templates/realtime.service
[Unit]
Description=nself Realtime Socket.io Server
After=network.target postgresql.service redis.service
Wants=postgresql.service redis.service

[Service]
Type=simple
User=nself
Group=nself
WorkingDirectory=/opt/nself-plugins/plugins/realtime/ts
Environment=NODE_ENV=production
EnvironmentFile=/opt/nself-plugins/plugins/realtime/.env
ExecStart=/usr/bin/node dist/server.js
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal
SyslogIdentifier=nself-realtime

# Security
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/opt/nself-plugins/plugins/realtime

# Resource limits
LimitNOFILE=65536
LimitNPROC=4096

[Install]
WantedBy=multi-user.target
